# build/resourced/cracklord-resourced.service
[Unit]
Description=Cracklord Resource Daemon
After=network.target

[Service]
ExecStart=/usr/bin/cracklord-resourced --conf="/etc/cracklord/resourced.conf" &
ExecReload=/bin/kill -SIGINT $MAINPID
KillMode=process
Restart=on-failure
Type=simple

[Install]
Alias=cracklord-resourced.service
