# build/resourced/cracklord-resourced.conf
description   "Cracklord Resource Server"
author        "Lucas Morris & Michael McAtee"

start on (net-device-up and local-filesystems and runlevel [2345])
stop  on runlevel [016]

respawn

chdir /var/cracklord/

setuid cracklord

exec /usr/bin/cracklord-resourced --conf="/etc/cracklord/resourced.conf" 
