# build/queued/cracklord-queued.conf
description    "Cracklord Queue Server"
author         "Lucas Morris & Michael McAtee"

start on (net-device-up and local-filesystems and runlevel [2345])
stop  on runlevel [016]

respawn

exec /usr/bin/cracklord-queued --conf="/etc/cracklord/queued.conf" --webroot="/var/cracklord/www/" --key="/etc/cracklord/ssl/queued.key" --cert="/etc/cracklord/ssl/queued.crt"
